# deploy/channel-gateway.service
[Unit]
Description=Channel Gateway (in-house channels -> browser, HLS)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust User and WorkingDirectory to match your server.
User=CHANGE_ME_USER
WorkingDirectory=/opt/channel-gateway
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=3
# Don't spam the journal if it crash-loops on a config error.
StartLimitIntervalSec=60
StartLimitBurst=10

[Install]
WantedBy=multi-user.target
